# bin/work/vlan_conf
#!/bin/bash

### BEGIN INIT INFO
# Provides:          Mehdi GHESH VLAN
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Set up VLAN for Mehdi GHESH's computer
# Description:       This script automatically set up the network
#                    settings
### END INIT INFO

# The original file name is : mgh_vlan_conf

# This script will auto link VLAN configured with the switch to
# be able to work either on a private network with the Box and
# also to have access to the internet

#############################################################################
# Global Knowledge
#############################################################################

# Global variables
status_file=/tmp/vlan_conf_stat

fail_unless_root() {
	me=`whoami`
	if [ $me != "root" ]
	then
		echo "The script must be run with root privileges"
		exit 1
	fi
}

# Commande qui permet d'afficher la configuration des interfaces réseaux
# ip a
# Commande qui permet d'afficher les routes des interfaces réseaux
# ip r


# Commande à utiliser lorsqu'on se connecte à un switch
# cela permet ensuite de pouvoir la configurer depuis l'interface graphique sur :80
# Commande qui assigne la plage d'adresse à l'interface
# sudo ip addr add 192.168.2.100/24 dev enp0s31f6

# Il est possible  que le network manager retire l'adresse précédemment et casse donc la connexion entre l'ordi et le switch
# pour parer cette événetualité, penser à arrêter le service
# sudo /etc/init.d/network-manager stop
# ping 192.168.2.10 # test de l'ajout

setup_new_vlan() {
	local name="$1"
	local ip="$2"
	local id="$3"
	# Ajout d'un VLAN, d'identifiant 1 et nommé sw0 à l'interface ethernet
	sudo ip link add link enp0s31f6 name ${name} type vlan id $id
	# Commande qui assigne la plage d'adresse au VLAN sw0
	sudo ip addr add ${ip} dev ${name}
	# Il faut maintenant activer les nouvelles modifications
	sudo ip link set ${name} up
}

setup_switch_vlan() {
	# Ajout d'un VLAN, d'identifiant 1 et nommé sw0 à l'interface ethernet
	sudo ip link add link enp0s31f6 name sw0 type vlan id 1
	# Commande qui assigne la plage d'adresse au VLAN sw0
	sudo ip addr add 192.168.2.100/24 dev sw0
	# Il faut maintenant activer les nouvelles modifications
	sudo ip link set sw0 up
}

unsetup_switch_vlan() {
	sudo ip link delete sw0
}

# sudo /etc/init.d/network-manager start # Relancement du network manager

setup_box_vlan() {
	# Ajout d'un VLAN, d'identifiant 2 et nommé box0 à l'interface ethernet
	sudo ip link add link enp0s31f6 name box0 type vlan id 2
	# Commande qui assigne la plage d'adresse au VLAN box0
	sudo ip addr add 192.168.1.100/24 dev box0
	# Il faut maintenant activer les nouvelles modifications
	sudo ip link set box0 up
}

unsetup_box_vlan() {
	sudo ip link delete box0
}

setup_efixo_dns() {
	# The file might have been attributs changed with : sudo chattr [+-]i /etc/resolv.conf
	# https://wiki.archlinux.fr/resolv.conf
	sudo sed -i -r 's/^nameserver 127.0.0.1/nameserver 192.168.100.254/' /etc/resolv.conf
	sudo dhclient
}

unsetup_efixo_dns() {
	true
}


case "$1" in
	start)
		setup_new_vlan management 192.168.50.100/24 50
		setup_new_vlan nb6vac 192.168.60.100/24 60
		setup_new_vlan gen8SDK 192.168.61.100/24 61
		setup_new_vlan gen8BSP 192.168.62.100/24 62

		# setup_efixo_dns
		echo "VLAN correctly configured" >>$status_file
		;;

	stop)
		unsetup_box_vlan
		unsetup_switch_vlan
		unsetup_efixo_dns
		echo "VLAN correctly stopped" >>$status_file
		;;

	restart)
		$0 stop
		$0 start
		;;

	force-reload)
		$0 restart
		;;

	status)
		cat $status_file
		;;

	*)
		echo "Usage: service $0 {start|stop|restart|status}"
		exit 1
		;;
esac

exit 0
